# config/hdcpd.service.in
[Unit]
Description=hdcp daemon

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_BINDIR@/hdcpd

[Install]
WantedBy=multi-user.target
